# qne-node.service
[Unit]
Description=QNE Node QUIC Server
After=network.target

[Service]
Type=simple
User=qne
ExecStart=/usr/local/bin/qne-node
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
